# deploy/systemd/healthcheck.service
[Unit]
Description=Stekkies bot health check (DeepSeek credit + Stekkies login)
After=browser-host.service
Requires=browser-host.service

[Service]
Type=oneshot
User=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=DISPLAY=__DISPLAY__
Environment=PATH=__APP_HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=__APP_DIR__/state/agent.env
ExecStart=__APP_HOME__/.local/bin/uv run python -m src.healthcheck
