# deploy/systemd/orchestrator.service
[Unit]
Description=Stekkies orchestrator (Gmail trigger -> extract -> apply agent)
After=browser-host.service
Requires=browser-host.service

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=DISPLAY=__DISPLAY__
Environment=PATH=__APP_HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
# Secrets (OPENROUTER_API_KEY, optional APPLY_MODEL overrides) live here:
EnvironmentFile=__APP_DIR__/state/agent.env
ExecStart=__APP_HOME__/.local/bin/uv run python -m src.orchestrator
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
