# deploy/systemd/poller.service
[Unit]
Description=Active listing poller (watch source sites -> apply agent)
After=browser-host.service
Requires=browser-host.service

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=DISPLAY=__DISPLAY__
Environment=PATH=__APP_HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
# Same secrets/profile as the orchestrator (DEEPSEEK_API_KEY, APPLICANT_*,
# NOTIFY_TO, plus optional POLL_* / POLL_JUDGE_MODEL overrides):
EnvironmentFile=__APP_DIR__/state/agent.env
ExecStart=__APP_HOME__/.local/bin/uv run python -m src.poller.watcher
Restart=always
RestartSec=10
# A manual `systemctl stop` (SIGTERM) is a clean pause, not a failure.
SuccessExitStatus=143 SIGTERM

[Install]
WantedBy=multi-user.target
