# deploy/systemd/self-improvement-worker.service
[Unit]
Description=Drain durable Stekkies self-improvement jobs serially
After=network-online.target litellm-proxy.service browser-host.service
Wants=network-online.target litellm-proxy.service browser-host.service

[Service]
Type=oneshot
User=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=DISPLAY=:99
Environment=PATH=__APP_HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=__APP_DIR__/state/agent.env
ExecStart=__APP_HOME__/.local/bin/uv run python -m src.self_improvement_worker
TimeoutStartSec=2h
